# supervisor/conf.d/telebot.conf
[program:telebot-polling]
command=php /anilibrary/artisan telebot:polling
autorestart=true
redirect_stderr=true
stdout_logfile=/anilibrary/storage/logs/telebot.log
stopwaitsecs=3600