# supervisor/conf.d/telebot.conf
[program:telebot-polling]
command=php artisan telebot:polling
autorestart=true
stdout_logfile=/anilibrary/storage/logs/telebot.log
stopwaitsecs=3600