# docker/supervisor/conf.d/horizon.conf
[program:horizon]
process_name=%(program_name)s
command=php artisan horizon
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/anilibrary/storage/logs/horizon.log
stopwaitsecs=3600